# systemd-files/webkpasswd.socket
[Unit]
Description=gunicorn webkpasswd socket

[Socket]
ListenStream=/run/webkpasswd/socket
User=nginx

[Install]
WantedBy=sockets.target
